# scripts/mac-mini/launchd/com.mana.telegram-ollama-bot.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.mana.telegram-ollama-bot</string>
    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/node</string>
        <string>/Users/mana/projects/managarten/services/telegram-ollama-bot/dist/main.js</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/mana/projects/managarten/services/telegram-ollama-bot</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>TELEGRAM_BOT_TOKEN</key>
        <string>8559479868:AAHF3g7vYLs0eOvDLh7hFVnIB-V8CKehUOM</string>
        <key>OLLAMA_URL</key>
        <string>http://localhost:11434</string>
        <key>OLLAMA_MODEL</key>
        <string>gemma3:4b</string>
        <key>PORT</key>
        <string>3301</string>
    </dict>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/mana/Library/Logs/telegram-ollama-bot.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/mana/Library/Logs/telegram-ollama-bot.log</string>
</dict>
</plist>
